# .ops/docker/config/crontab
* * * * * artisan schedule:run >> /var/www/html/src/Infrastructure/Boilerplate/Laravel/storage/logs/laravel.scheduler.log 2>&1
