# .ops/docker/config/dcron.conf
[program:dcron]
command=crond -f -S
autostart=true
autorestart=true
priority=15
stdout_events_enabled=true
stderr_events_enabled=true
